# openwrt部署
STEP1:
vi /etc/init.d/cf_ddns


STEP2:
#!/bin/sh /etc/rc.common
START=99
USE_PROCD=1

start_service() {
    procd_open_instance
    # 修改：这里填写你脚本的绝对路径
    procd_set_param command /bin/bash /root/ddns/cloudflare_DDNS_hao.sh
    # 自动重启
    procd_set_param respawn
    # 即使报错输出也会重定向到系统日志
    procd_set_param stderr 1
    procd_set_param stdout 1
    procd_close_instance
}

STEP3:
# 赋权限
chmod +x /etc/init.d/cf_ddns
# 启用开机自启
/etc/init.d/cf_ddns enable
# 立即启动
/etc/init.d/cf_ddns start
